# automation/com.codex.nba.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex.nba.daily</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/ykonsta/Desktop/CodeX_Models/.venv/bin/python3</string>
      <string>/Users/ykonsta/Desktop/CodeX_Models/scripts/run_daily.py</string>
      <string>--days-ahead</string>
      <string>7</string>
      <string>--odds-provider</string>
      <string>auto</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ykonsta/Desktop/CodeX_Models</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ykonsta/Desktop/CodeX_Models/logs/launchd_daily.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ykonsta/Desktop/CodeX_Models/logs/launchd_daily.err.log</string>

    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
